# etc/nova-agent.debian
#!/bin/sh

# LSB style init header:
#
### BEGIN INIT INFO
# Provides: Nova-Agent
# Required-Start: $remote_fs $syslog xe-linux-distribution
# Required-Stop: $remote_fs $syslog
# Default-Start: 2 3 4 5
# Default-Stop: 0 1 6
# Short-Description: Start nova-agent at boot time
# Description: nova-agent is a guest agent for OpenStack nova.
### END INIT INFO

# Source function library.
. /lib/lsb/init-functions

PROG=nova-agent
exec=/usr/bin/$PROG
pid_file=/var/run/nova-agent.pid
LOCK=/var/lock/subsys/$PROG

do_start() {
    if pidofproc -p $pid_file; then
        return 0
    fi
    export PATH="/sbin:/bin:/usr/sbin:/usr/bin:/usr/games:/usr/local/sbin:/usr/local/bin"
    start-stop-daemon --background --pidfile $pid_file --start --exec $exec -- -o /var/log/nova-agent.log -l info
    RETVAL=$?
    echo "Starting $PROG... "
    if [ $RETVAL -eq 0 ]; then
        [ -d /var/lock/subsys ] || mkdir -p /var/lock/subsys
        touch $LOCK
    fi
    return $RETVAL
}

do_stop() {
    start-stop-daemon --pidfile $pid_file --stop
    RETVAL=$?
    echo "Stopping $PROG... "
    if [ $RETVAL -eq 0 ]; then
        rm -f $LOCK
    fi
    return $RETVAL
}

case "$1" in
  start)
    do_start
    ;;
  stop)
    do_stop
    ;;
  restart)
    do_stop
    do_start
    ;;
  *)
    echo "Usage: $SCRIPTNAME {start|stop|restart}" >&2
    exit 3
    ;;
esac
