# etc/nova-agent.gentoo
#!/sbin/runscript

depend() {
    before net
    after xe-daemon
}

start() {
    command="/usr/bin/nova-agent"
    pidfile="/var/run/nova-agent.pid"

    ebegin "Starting nova-agent"
    start-stop-daemon --background --start --exec $command \
        --pidfile $pidfile -- -p $pidfile -o /var/log/nova-agent.log
    eend $?
}

stop() {
    command="/usr/bin/nova-agent"
    pidfile="/var/run/nova-agent.pid"
    ebegin "Stopping nova-agent"
    start-stop-daemon --stop --exec ${command} \
        --pidfile ${pidfile} --quiet
    eend $?
}
