# etc/nova-agent.upstart
description "Nova Agent upstart script"
author          "Daniel Wallace"

start on networking or runlevel [2345]
stop on runlevel [!2345]

respawn
respawn limit 10 5
umask 022

exec /usr/bin/nova-agent -o /var/log/nova-agent.log
