# etc/novaagent.service
[Unit]
Description=Nova Agent for xenstore

[Service]
ExecStart=/usr/bin/nova-agent

[Install]
WantedBy=network.target
