# etc/nova-agent.alpine
#!/sbin/openrc-run
#
# nova-agent
#


name="nova-agent"

description="Agent for setting up clean servers on Xen"
description_reload="Restart nova-agent"

command="/usr/bin/$name"

pidfile="/var/run/nova-agent.pid"
lockfile="/var/lock/subsys/$name"
logfile="/var/log/nova-agent.log"

loglevel="info"

extra_started_commands="reload"

start_pre() {
  [ -x $command ] || exit 5
}

start() {
  ebegin $"Starting $name"
  start-stop-daemon \
    --signal HUP \
    --exec "$command" \
    --pidfile "$pidfile" \
    -- \
      -o "$logfile" \
      -l "$loglevel"
  eend $?
}

stop() {
  ebegin $"Stopping $name"
  start-stop-daemon \
    --stop \
    --exec "$command" \
    --pidfile "$pidfile" \
    --quiet
  eend $?
}

restart() {
  stop
  start
}

reload() {
  restart
}

force_reload() {
  restart
}
