# etc/nova-agent.upstart
description "Nova Agent upstart script"
author      "Rackspace"

start on networking or runlevel [2345]
stop on runlevel [!2345]

respawn
respawn limit 10 5
expect stop
umask 022

exec /usr/bin/nova-agent --no-fork -o /var/log/nova-agent.log -l info
